# CONFIG/wlan0_ap-heal.timer
[Unit]
Description=Periodic check & heal for wlan0_ap

[Timer]
OnBootSec=30
OnUnitActiveSec=60

[Install]
WantedBy=timers.target

